# docker/rc.local
#!/bin/sh
/opt/dell/srvadmin/sbin/srvadmin-services.sh restart
